# packaging/systemd/kimi-work-update.timer
[Unit]
Description=Weekly update check for Kimi Work (unofficial Linux build)

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
